# lesson-1/crontab
* * * * * php /mnt/c/Users/Василий/Documents/practicecode/php-cli/lesson-1/script.php >> /mnt/c/Users/Василий/Documents/log.txt